# eyejoker/native-claude/ops/claude-native-cohort-verifier.service
[Unit]
Description=Verify new Claude Code and Agent SDK cohorts in isolation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ejclaw/projects/ejclaw-native/eyejoker/native-claude
EnvironmentFile=-/home/ejclaw/.config/claude-native/env
UnsetEnvironment=DISCORD_BOT_TOKEN
Environment=PATH=/home/ejclaw/.bun/bin:/home/ejclaw/.hermes/node/bin:/home/ejclaw/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=npm_config_cache=/home/ejclaw/.local/state/claude-native/npm-cache
ExecStart=/home/ejclaw/.bun/bin/bun run src/cohort-verifier.ts
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ejclaw/.local/state/claude-native
RestrictSUIDSGID=true
TimeoutStartSec=45min
